# utils/test/reporting/docker/supervisor.conf
[supervisord]
nodaemon = true

[program:reporting_tornado]
user = root
directory = /home/opnfv/utils/test/reporting/api
command = python server.py --port=800%(process_num)d
process_name=%(program_name)s%(process_num)d
numprocs=4
numprocs_start=1
autorestart = true

[program:reporting_nginx]
user = root
command = service nginx restart
autorestart = true
