# utils/test/reporting/docker/supervisor.conf
[supervisord]
nodaemon = true

[program:tornado]
user = root
directory = /home/opnfv/releng/utils/test/reporting/api/api
command = python server.py --port=800%(process_num)d
process_name=%(program_name)s%(process_num)d
numprocs=4
numprocs_start=1

[program:nginx]
user = root
command = service nginx restart

[program:configuration]
user = root
directory = /home/opnfv/releng/utils/test/reporting/pages
command = bash config.sh
